# NSWI106_linux_admin/05/00-systemd-user.d/vde@.service
[Unit]
Description=vde switch %i startup service


[Service]
Type=forking
ExecStart=/usr/bin/vde_switch -daemon -sock /home/tomisz/vde/sw%i.sock

[Install]
WantedBy=multi-user.target
RequiredBy=gw.service ns%i.service
